# files/ami/etc/local.d/00_cloud_init.start
# cloud-init init
hostname `curl -s http://169.254.169.254/2008-02-01/meta-data/hostname`
if [ ! -d /root/.ssh ]; then
	mkdir -p /root/.ssh
	chmod 0700 /root/.ssh
fi
wget -q -O /root/.ssh/authorized_keys http://169.254.169.254/2008-02-01/meta-data/public-keys/0/openssh-key

ifconfig eth0 mtu 1454
