# files/ami/etc/local.d/01_firstboot.start
features=`grep flags /proc/cpuinfo |uniq|cut -c 10-`
cpus=`grep -c processor /proc/cpuinfo`
echo "MAKEOPTS=\"-j${cpus}\"" >> /etc/portage/make.conf
echo "USE=\"\${USE} ${features}\"" >> /etc/portage/make.conf
echo 'USE="${USE} mysql sqlite python conntrack augeas virt-network -monitor"' >> /etc/portage/make.conf

export NOCOLOR="yes"
/usr/bin/emerge-webrsync
time emerge -q eix && eix-update || exit

chmod -x /etc/local.d/01_firstboot.start
